# com.trusttunnel.gui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trusttunnel.gui</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/trusttunnel-macos/.venv/bin/python3</string>
        <string>-m</string>
        <string>src.app</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/trusttunnel_client</string>
    </dict>
</dict>
</plist>
